# systemd/playlog-backend.service
[Unit]
Description=playlog backend
Wants=network-online.target
After=network-online.target

[Service]
User=playlog
Type=simple
WorkingDirectory=/var/lib/playlog/playlog
ExecStart=/var/lib/playlog/playlog/playlog -vl3002
Restart=on-failure
RestartSec=10
RestartSteps=60
RestartMaxDelaySec=900

[Install]
WantedBy=multi-user.target
